# server/systemd/myreactor-gunicorn.socket
[Unit]
Description=gunicorn socket for the myreactor project

[Socket]
ListenStream=/run/myreactor-gunicorn.sock

[Install]
WantedBy=sockets.target